# ops/systemd/peptaide-supabase.service
[Unit]
Description=Peptaide local Supabase stack (supabase CLI)
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=oneshot
WorkingDirectory=/data/projects/peptaide
ExecStart=/home/ubuntu/.local/bin/supabase start
ExecStop=/home/ubuntu/.local/bin/supabase stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
